# deploy/dockers/flask/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:flask]
command=/bin/bash -c "python /src/run.py"
